# upstart/semantic.conf
# /etc/init/semantic.conf

description "SMART semantic services"

start on startup
stop on shutdown
respawn
respawn limit 20 5

script
   cd /usr/smart/semantic
   node semantic
end script
